# templates/default/templar-1.conf.erb
start on starting templar
stop on stopping templar
respawn

exec templar -redis='<%= @redis %>' -listen='<%= @listen %>' -debug='<%= @debug %>' > <%= @log_location %>
